# chef/cookbooks/jira/templates/default/jira.erb
#!/bin/sh
#
# PROVIDE: jira
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="jira"
rcvar=`set_rcvar`

load_rc_config $name

: ${jira_enable="NO"}

pidfile='<%= @node[:jira_install_path] %>/jira.pid'
procname='java'

export JAVA_HOME=/usr/local/diablo-jre1.6.0
export CATALINA_PID=$pidfile

start_cmd="su -m <%= @node[:jira_run_user] %> -c 'sh -c \"<%= @node[:jira_install_path] %>/bin/startup.sh\"'"
stop_cmd="<%= @node[:jira_install_path] %>/bin/shutdown.sh"

run_rc_command "$1"
